# .config/udev/rules.d/99-wifi-power-save.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wl*", ATTR{power/control}=="auto", ATTR{power/control}="on"
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wl*", RUN+="/usr/bin/iw dev $name set power_save off"
